# dev-python/python-matter-server/files/matter-server.service
# /etc/systemd/system/matter-server.service
[Unit]
Description=Python Matter Server (for Home Assistant Matter integration)
Documentation=https://github.com/home-assistant-libs/python-matter-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=matter-server
Group=matter-server
WorkingDirectory=/var/lib/matter-server
ExecStart=/usr/bin/matter-server --storage-path /var/lib/matter-server/storage --paa-root-cert-dir /var/lib/matter-server/credentials
Restart=on-failure
RestartSec=10
# Optional: adjust log level if needed
# Environment=MATTER_SERVER_LOG_LEVEL=info

# Security hardening
NoNewPrivileges=yes
#ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
